# bot/deploy/kidgo-bot.service
[Unit]
Description=KidGo multi-platform bot (Telegram/WhatsApp/Signal)
# The bot calls the local API, so start it after that's up.
After=network.target kidgo-api.service
Wants=kidgo-api.service

[Service]
Type=simple
WorkingDirectory=/opt/kidgo-bot
EnvironmentFile=/opt/kidgo-bot/.env
# Cap the heap so the bot can't starve the API/Postgres on the ~512MB host.
# grammY is lightweight (long polling, no browser); 128MB is plenty.
ExecStart=/usr/bin/node --max-old-space-size=128 src/index.js
# Always restart so a hard crash never leaves the bot silently dead.
Restart=always
RestartSec=5
MemoryMax=220M

[Install]
WantedBy=multi-user.target
